# etc/profile.d/vx-desktop-settings.sh
# si no se establece esta variable, programas como Evince siguen mostrando
# las scrollbars de Unity
export LIBOVERLAY_SCROLLBAR=0
